# runit/metadata-collector/run
#!/bin/sh
exec python3 -u -m metadata_collector